# nginx/supervisor-app.conf
[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /code/uwsgi.ini  --py-autoreload 1

[program:nginx-app]
command = /usr/sbin/nginx